# devops/linux/openclaw-workspace-backup.timer
[Unit]
Description=OpenClaw Workspace Backup Timer

[Timer]
OnBootSec=10min
# FIXED: Reduced from 4h to 2h to reduce potential data loss window
# Source: unisone/openclaw-config production hardening recommendation
OnUnitActiveSec=2h
AccuracySec=5min

[Install]
WantedBy=timers.target
