# devops/linux/openclaw-health-check.timer
[Unit]
Description=OpenClaw Health Check Timer

[Timer]
OnBootSec=5min
# FIXED: Reduced from 30min to 5min for faster crash detection
# Source: unisone/openclaw-config production hardening research
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
